# systemd/echo-listen.service
[Unit]
Description=Echo social media agent - listening mode
After=network-online.target

[Service]
Type=simple
User=<YOUR_USER>
EnvironmentFile=/srv/ai/echo/configs/echo.env
ExecStart=/usr/bin/python3 /srv/ai/echo/echo.py listen --quiet
Restart=on-failure
RestartSec=60
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/srv/ai/echo

[Install]
WantedBy=multi-user.target
